# scripts/backup/systemd/homegrown-backup.timer
[Unit]
Description=Homegrown Academy — daily PostgreSQL backup timer
Requires=homegrown-backup.service

[Timer]
# Run daily at 02:30 UTC, randomised by up to 10 minutes to avoid thundering herd.
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
