# scripts/backup/systemd/homegrown-restore-drill.timer
[Unit]
Description=Homegrown Academy — weekly restore drill timer
Requires=homegrown-restore-drill.service

[Timer]
# Run every Sunday at 04:00 UTC.
OnCalendar=Sun *-*-* 04:00:00 UTC
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
